# deploy/unichat-scheduler.service
[Unit]
Description=uni-chat routines scheduler (APScheduler)
After=network-online.target mongod.service unichat.service
Wants=network-online.target

[Service]
Type=simple
User=unichat
Group=unichat
WorkingDirectory=/opt/unichat/scheduler
EnvironmentFile=/etc/unichat/scheduler.env
ExecStart=/opt/unichat/scheduler/.venv-uv/bin/python -m scheduler.main
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60
LimitNOFILE=65535

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
