# deploy/unichat-bot.service
[Unit]
Description=uni-chat Telegram bot (aiogram polling)
After=network-online.target mongod.service unichat.service
Wants=network-online.target

[Service]
Type=simple
User=unichat
Group=unichat
WorkingDirectory=/opt/unichat/bot
EnvironmentFile=/etc/unichat/bot.env
ExecStart=/opt/unichat/bot/.venv-uv/bin/python -m bot.main
Restart=on-failure
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
